# systemd/garden-sensors-gatherer.service
[Unit]
Description=Garden sensors data gatherer

[Service]
User=root
Type=simple
ExecStart=/home/masneyb/data/garden-sensors/bin/gather-data.sh
