# systemd/garden-sensors-create-graphs.service
[Unit]
Description=Garden sensors graph creator

[Service]
User=masneyb
Type=simple
ExecStart=/home/masneyb/data/garden-sensors/bin/create-graphs.sh
