# systemd/weather-station-create-graphs.timer
[Unit]
Description=Garden sensors graph creator timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
