# systemd/weather-station-gatherer.service
[Unit]
Description=Garden sensors data gatherer

[Service]
User=root
Type=simple
ExecStart=/home/masneyb/data/weather-station/bin/gather-data.sh
