# pi-yadl/systemd/pi-yadl-gatherer.timer
[Unit]
Description=pi-yadl data gatherer timer

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=multi-user.target
