# pi-yadl/systemd/pi-yadl-create-graphs.timer
[Unit]
Description=pi-yadl graph creator timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
